# src/s6-based-images/roles/s6-server/files/services.d/sshd/run
#!/usr/bin/execlineb -P

# exec /usr/bin/one_time.sh

/usr/sbin/sshd -D
