# src/wiremock/files/services.d/wiremock/run
#!/usr/bin/with-contenv sh

exec /usr/bin/wiremock.sh
