# src/s6-based-images/roles/vault-server/files/s6-overlay/s6-rc.d/unseal-vault/up
/command/with-contenv /usr/local/bin/unseal.py
