# src/spark/hadoop/files/services.d/hadoop/run
#!/usr/bin/execlineb -P

s6-setuidgid app

/home/app/bin/master.sh
