# src/spark/files/services.d/master/run
#!/usr/bin/with-contenv sh

exec /usr/local/spark/bin/master.sh