# src/s6-based-images/roles/vault-server/files/s6-overlay/s6-rc.d/vault/run
#!/command/with-contenv sh

FILE=/data/raft-vault/vault.db
if [ ! -f "$FILE" ]; then
  cp -r /backup/* /data/raft-vault/
fi

vault server -log-level=trace -config=/configs/config-vault.hcl