# src/s6-based-images/roles/vault-server/files/s6-overlay/s6-rc.d/traefik/run
#!/command/with-contenv sh

traefik --configfile /configs/traefik.yaml
