# src/serv-virtualization/services/files/services.d/webhook-transformer/run
#!/usr/bin/with-contenv sh

exec /usr/bin/webhook-transformer.sh
