# src/spark/files/services.d/worker/run
#!/usr/bin/with-contenv sh

exec /usr/local/spark/bin/worker.sh