# src/serv-virtualization/services/files/services.d/service-a/run
#!/usr/bin/with-contenv sh

exec /home/app/service-a.sh